# configs/php-fpm-supervisor.conf
[program:php-fpm]
command=php-fpm -F
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/wp_app/php-fpm.log